# systemd/tunasync-snapshot-gc.service
[Unit]
Description=Delete garbage subvolumes generated by tunasync
Requires = network.target
After = network.target

[Service]
Type=simple
ExecStart=/home/tuna/.virtualenvs/tunasync/bin/python -u /home/tuna/tunasync/tunasync_snapshot_gc.py -c /etc/tunasync.ini

[Install]
WantedBy = multi-user.target
